# deploy/gbgolf.service
[Unit]
Description=GB Golf Optimizer (Gunicorn)
After=network.target

[Service]
User=deploy
Group=www-data
WorkingDirectory=/opt/GBGolfOptimizer
Environment="PATH=/opt/GBGolfOptimizer/.venv/bin"
Environment="SCRIPT_NAME=/golf"
ExecStart=/opt/GBGolfOptimizer/.venv/bin/gunicorn \
    --workers 2 \
    --bind unix:/opt/GBGolfOptimizer/gbgolf.sock \
    --umask 007 \
    wsgi:app
Restart=on-failure

[Install]
WantedBy=multi-user.target
